# systemd/ipnameset.service
[Unit]
Description=Update ipnameset name based ipset

[Service]
ExecStart=/usr/local/bin/ipnameset-helper populate
